# tasks/ebpf_jit/overlay/etc/profile.d/home-root.sh
#! /bin/sh -

export HOME=/root
cd "${HOME}" || exit
